# launchd/com.selene.dev-process-batch.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.selene.dev-process-batch</string>
    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <string>-c</string>
        <string>cd /Users/chaseeasterling/selene-n8n &amp;&amp; ./scripts/dev-process-batch.sh 15</string>
    </array>
    <key>StartInterval</key>
    <integer>3600</integer>
    <key>EnvironmentVariables</key>
    <dict>
        <key>SELENE_ENV</key>
        <string>development</string>
        <key>PATH</key>
        <string>/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin</string>
    </dict>
    <key>StandardOutPath</key>
    <string>/Users/chaseeasterling/selene-data-dev/logs/batch-process.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/chaseeasterling/selene-data-dev/logs/batch-process-error.log</string>
    <key>RunAtLoad</key>
    <false/>
</dict>
</plist>
